# layers/meta-olympus/recipes-connectivity/olympus-ap/files/olympus-ap.service
[Unit]
Description=Olympus AP — rover WiFi hotspot (field mode)
Documentation=https://github.com/Alonso11/olympus-hlc-rpi5
Conflicts=wpa_supplicant@wlan0.service
After=network.target

[Service]
Type=forking
PIDFile=/run/hostapd.pid
RemainAfterExit=no

# Detener cliente WiFi si está activo
ExecStartPre=-/bin/systemctl stop wpa_supplicant@wlan0.service

# Configurar wlan0 con IP fija del rover
ExecStartPre=/sbin/ip link set wlan0 down
ExecStartPre=/sbin/ip addr flush dev wlan0
ExecStartPre=/sbin/ip link set wlan0 up
ExecStartPre=/sbin/ip addr add 192.168.100.1/24 dev wlan0

# Lanzar hostapd en background
ExecStart=/usr/sbin/hostapd -B -P /run/hostapd.pid /etc/hostapd/hostapd.conf

# Lanzar dnsmasq para DHCP
ExecStartPost=/usr/sbin/dnsmasq \
    --conf-file=/etc/dnsmasq-ap.conf \
    --pid-file=/run/dnsmasq-ap.pid

ExecStop=-/bin/sh -c 'kill $$(cat /run/hostapd.pid) 2>/dev/null; true'
ExecStop=-/bin/sh -c 'kill $$(cat /run/dnsmasq-ap.pid) 2>/dev/null; true'
ExecStop=/sbin/ip addr flush dev wlan0

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
